# Docker/nginx/crontab
35 3 * * * /usr/bin/sh /certbot/certbot_renewal.sh > /dev/stdout
*/20 * * * * /usr/bin/sh /duckdns/dns_update.sh > /dev/stdout
