# Docker/mariadb/crontab
30 2 * * * BASH_ENV=/etc/profile /usr/bin/sh /backup.sh > /backup
